# crontabs/kitweather_joblist.crontab
#-----------------------------------------------------------------
# Shell variable for cron
SHELL=/bin/bash
#
# PATH variable for cron
PATH=/home/iconeps/miniconda3/bin:/usr/local/bin:/usr/bin:/usr/local/sbin:/usr/sbin
#
#M   S     T M W   Befehl
#-----------------------------------------------------------------
0    0             * * *   python /progs/operational/w2w_ensembleplots/cronscripts/delete_old_forecast_data.py > /data/logs/scripts/delete_old_forecast_data.log 2> /data/logs/errors/delete_old_forecast_data.log
31   5,17          * * *   python /progs/operational/w2w_ensembleplots/cronscripts/download_icon-eu-eps.py > /data/logs/scripts/download_icon-eu-eps.log 2> /data/logs/errors/download_icon-eu-eps.log
50   10,22         * * *   python /progs/operational/w2w_ensembleplots/cronscripts/download_icon-eu-eps.py > /data/logs/scripts/download_icon-eu-eps.log 2> /data/logs/errors/download_icon-eu-eps.log
26   5,17          * * *   python /progs/operational/w2w_ensembleplots/cronscripts/download_icon-global-eps.py > /data/logs/scripts/download_icon-global-eps.log 2> /data/logs/errors/download_icon-global-eps.log
49   4,16          * * *   python /progs/operational/w2w_ensembleplots/cronscripts/download_icon-eu-det.py > /data/logs/scripts/download_icon-eu-det.log 2> /data/logs/errors/download_icon-eu-det.log
43   10,22         * * *   python /progs/operational/w2w_ensembleplots/cronscripts/download_icon-eu-det.py > /data/logs/scripts/download_icon-eu-det.log 2> /data/logs/errors/download_icon-eu-det.log
40   4,16          * * *   python /progs/operational/w2w_ensembleplots/cronscripts/download_icon-global-det.py > /data/logs/scripts/download_icon-global-det.log 2> /data/logs/errors/download_icon-global-det.log
15   7,19          * * *   python /progs/operational/w2w_ensembleplots/cronscripts/meteogram_boxplot_w2w_cities_latest.py > /data/logs/scripts/meteogram_boxplot_w2w_cities_latest_0_12UTC.log 2> /data/logs/errors/meteogram_boxplot_w2w_cities_latest_0_12UTC.log
30   11,23         * * *   python /progs/operational/w2w_ensembleplots/cronscripts/meteogram_boxplot_w2w_cities_latest.py > /data/logs/scripts/meteogram_boxplot_w2w_cities_latest_6_18UTC.log 2> /data/logs/errors/meteogram_boxplot_w2w_cities_latest_6_18UTC.log
0    5,17          * * *   python /progs/operational/w2w_ensembleplots/cronscripts/det_overview_map_precrate_mslp.py > /data/logs/scripts/det_overview_map_precrate_mslp.log 2> /data/logs/errors/det_overview_map_precrate_mslp.log
40   5,17          * * *   python /progs/operational/w2w_ensembleplots/cronscripts/det_overview_map_t850_z500.py > /data/logs/scripts/det_overview_map_t850_z500.log 2> /data/logs/errors/det_overview_map_t850_z500.log
20   6,18          * * *   python /progs/operational/w2w_ensembleplots/cronscripts/det_overview_map_the850_mslp.py > /data/logs/scripts/det_overview_map_the850_mslp.log 2> /data/logs/errors/det_overview_map_the850_mslp.log
0    7,19          * * *   python /progs/operational/w2w_ensembleplots/cronscripts/det_overview_map_wind300_z300.py > /data/logs/scripts/det_overview_map_wind300_z300.log 2> /data/logs/errors/det_overview_map_wind300_z300.log
#-----------------------------------------------------------------